# csv-analyzer-mcpo.service
[Unit]
Description=CSV Analyzer MCPO Server
After=network.target

[Service]
Type=simple
User=mcpo
Group=mcpo
WorkingDirectory=/opt/csv-analyzer-mcpo
EnvironmentFile=/opt/csv-analyzer-mcpo/.env
ExecStart=/opt/csv-analyzer-mcpo/venv/bin/mcpo --host 0.0.0.0 --port 8091 --cors-allow-origins "*" -- /opt/csv-analyzer-mcpo/venv/bin/python /opt/csv-analyzer-mcpo/mcp_server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target